# systemd/kbd-backlight-on.service
[Unit]
Description=Keyboard backlight ON (AC only)

[Service]
Type=oneshot
ExecStart=/usr/local/lib/power-switch/kbd-backlight-on.sh