# systemd/kbd-backlight-battery.service
[Unit]
Description=Keyboard backlight OFF on battery

[Service]
Type=oneshot
ExecStart=/usr/local/lib/power-switch/kbd-backlight-battery.sh