# udev/98-power-backlight.rules
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/bin/systemctl start kbd-backlight-battery.service"