# udev/99-kbd-backlight.rules
ENV{ID_INPUT_KEYBOARD}=="1", ACTION=="change", RUN+="/bin/systemctl start kbd-backlight-on.service"